# js/proxy.conf
[program:chexie_net_proxy]
command=socat "TCP-LISTEN:80,fork,reuseaddr,su=nobody" "TCP6:chexie.net:443"
user=root
autostart=true
autorestart=true
redirect_stderr=true
